# lubuntu/.config/autostart/conky-left-autostart.desktop
[Desktop Entry]
Type=Application
Name=Conky Left autostart
Comment=Start Conky
Exec=conky -c ~/.conkyrc-test
