# grabbit-systemd/dot-config/systemd/user/homebrew.grabbit.service
[Unit]
Description=Run grabbit

[Service]
Type=oneshot
ExecStart=/home/linuxbrew/.linuxbrew/bin/grabbit grab
